# systemd/.config/systemd/user/tmux-save-sessions.service
[Unit]
Description=Tmux session saver

[Service]
Type=oneshot
ExecStart=zsh -c '$$XDG_DATA_HOME/tmux/save-sessions.sh'
